# supervisord.conf
[supervisord]
nodaemon=true

[program:X11]
command=Xvfb :0 -screen 0 1200x800x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc
autorestart=true

[program:novnc]
command=novnc_server --vnc localhost:5900 --listen 8080
autorestart=true

[program:socat]
command=socat tcp-listen:6000,reuseaddr,fork unix:/tmp/.X11-unix/X0
autorestart=true

[program:browsepy]
command=/usr/bin/browsepy 0.0.0.0 8000 --directory /root/ --upload /root/
stdout_events_enabled=true
stderr_events_enabled=true

[program:xfce]
command=startxfce4
autorestart=true

[program:pythonhttp]
command=/bin/bash -c "cd /ddkdocs/site && /usr/bin/python -m SimpleHTTPServer 7000"
autorestart=true
